# systemd/playbox.service
[Unit]
Description=playbox

[Service]
ExecStart=/home/pi/playbox/playbox.rb

[Install]
WantedBy=multi-user.target